# pg_backup/crontab
* * * * * /app/backup_and_send.sh >> /var/log/cron.log 2>&1
